# smssh/smssh.service
[Unit]
Description=ASMU — мониторинг SSH-атак
After=network.target sshd.service
Wants=sshd.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/smssh monitor
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target